# trading-ai/docs/systemd/ezra-research.service
[Unit]
Description=Ezra Trading AI (RESEARCH) — operating system supervisor
PartOf=ezra-trading.target
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/ezra-public/trading-ai

Environment=EZRAS_RUNTIME_ROOT=/opt/ezra-runtime
# Default safety: non-live. EnvironmentFile overlays may override when operator enables live stack.
Environment=NTE_EXECUTION_MODE=paper
Environment=NTE_LIVE_TRADING_ENABLED=false
Environment=NTE_PAPER_MODE=true
Environment=NTE_DRY_RUN=true
Environment=COINBASE_EXECUTION_ENABLED=false
Environment=COINBASE_ENABLED=false

EnvironmentFile=/opt/ezra-runtime/env/common.env
EnvironmentFile=-/opt/ezra-runtime/env/research.env
EnvironmentFile=-/opt/ezra-runtime/env/research-live.env
EnvironmentFile=-/opt/ezra-runtime/env/secrets.env

# Authoritative module resolution: prefer public tree (private may lag).
Environment=PYTHONPATH=/opt/ezra-public/trading-ai/src:/opt/ezra-private/trading-ai/src

ExecStart=/opt/ezra-venv/bin/python -m trading_ai.runtime daemon --role research --interval-sec 2.0 --skip-models
Restart=always
RestartSec=5

KillSignal=SIGINT
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target ezra-trading.target

